# thirdparty/includes-osx/Frameworks/CarbonCore.framework/Versions/Current/Resources/com.apple.coreservicesd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.coreservicesd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.CoreServices.coreservicesd</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/CoreServices/coreservicesd</string>
	</array>
    <key>HopefullyExitsLast</key>
    <true/>
	<key>EnableTransactions</key>
	<true/>
	<key>MachExceptionHandler</key>
	<string>com.apple.ReportCrash.DirectoryService</string>
</dict>
</plist>
